# poky/meta-st-openstlinux/recipes-bsp/tools/usbotg-gadget-config/97-ustotg.rules
ACTION=="add", SUBSYSTEM=="udc", RUN+="/sbin/stm32_usbotg_eth_config.sh start"
ACTION=="remove", SUBSYSTEM=="udc", RUN+="/sbin/stm32_usbotg_eth_config.sh stop"

